# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=batsignal -d 0 -c 5 -w 15 -W "battery low" -C "please bro don't let me die please just plug me in"
